# supervisord.conf
[supervisord]
nodaemon=true

[program:zeroview]
command=/usr/bin/python /opt/app/app.py
autorestart=true

[program:worker]
command=/usr/bin/python /opt/app/worker.py
autorestart=true
